# etc/crontabs/root
54 * * * *        /usr/sbin/libremap-agent
* * * * * /usr/sbin/reset_deaf_phys.sh
*/5 * * * * ( for file in /etc/alfred/* ; do [ -x $file ] && $file ; done )
